# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/etc/init.d/nginx start
